# package-manager.service.in
[Unit]
Description=Tizen Package Manager

[Service]
User=app_fw
Group=app_fw
ExecStart=@PREFIX@/bin/pkgmgr-server
