# package-manager.service.in
[Unit]
Description=Tizen Package Manager

[Service]
ExecStart=@PREFIX@/bin/pkgmgr-server
CapabilityBoundingSet=~CAP_MAC_ADMIN
CapabilityBoundingSet=~CAP_MAC_OVERRIDE
