# pkgmgr-server.service.in
[Unit]
Description=Package Manager daemon

[Service]
ExecStart=@PREFIX@/bin/pkgmgr-server
CapabilityBoundingSet=~CAP_MAC_ADMIN
CapabilityBoundingSet=~CAP_MAC_OVERRIDE
